# deploy/freebsd/freefsm
#!/bin/sh

# PROVIDE: freefsm
# REQUIRE: LOGIN postgresql
# KEYWORD: shutdown

. /etc/rc.subr

name="freefsm"
rcvar="freefsm_enable"

command="/usr/sbin/daemon"
procname="freefsm"
pidfile="/var/run/${name}/${name}.pid"
required_files="/usr/local/bin/freefsm"
freefsm_user="freefsm"

freefsm_config="/usr/local/etc/freefsm.conf"

load_rc_config $name

: ${freefsm_enable:="NO"}

start_precmd="${name}_prestart"

freefsm_prestart()
{
	if [ -f "${freefsm_config}" ]; then
		. "${freefsm_config}"
		export DATABASE_URL SESSION_SECRET PORT ENV STATIC_PATH SETUP_TOKEN
	fi

	if [ ! -d "/var/run/${name}" ]; then
		mkdir -p "/var/run/${name}"
	fi
	chown ${freefsm_user}:${freefsm_user} "/var/run/${name}"

	if [ ! -f /var/log/freefsm.log ]; then
		touch /var/log/freefsm.log
		chown ${freefsm_user}:${freefsm_user} /var/log/freefsm.log
		chmod 644 /var/log/freefsm.log
	fi

	command_args="-p ${pidfile} -u ${freefsm_user} -o /var/log/freefsm.log /usr/local/bin/freefsm"
}

run_rc_command "$1"
