# deploy/freebsd/freefsm
#!/bin/sh

# PROVIDE: freefsm
# REQUIRE: LOGIN postgresql
# KEYWORD: shutdown

. /etc/rc.subr

name="freefsm"
rcvar="freefsm_enable"
freefsm_user="freefsm"

start_cmd="freefsm_start"
status_cmd="freefsm_status"
stop_cmd="freefsm_stop"
stop_postcmd="freefsm_wait"
command="/usr/local/bin/freefsm"
daemon_command="/usr/sbin/daemon"
pidfile="/var/run/${name}.pid"

load_rc_config $name
: ${freefsm_enable:=no}

: ${freefsm_config:="/usr/local/etc/freefsm.conf"}

# set up dependant variables
procname="${command}"
required_files="${freefsm_config}"

freefsm_status()
{
  local pid
  pid=$(check_pidfile "${pidfile}" "${procname}")
  if [ -z "${pid}" ]
  then
    echo "FreeFSM is not running"
    return 1
  else
    echo "FreeFSM running, pid: ${pid}"
  fi
}

freefsm_start()
{
  echo "Starting freefsm:"
  ${daemon_command} -u "${freefsm_user}" -p "${pidfile}" -o /var/log/freefsm.log \
    ${command} \
    --config="${freefsm_config}"
}

freefsm_stop()
{
  echo "Stopping FreeFSM:"
  pid=$(check_pidfile "${pidfile}" "${procname}")
  if [ -z "${pid}" ]
  then
    echo "FreeFSM is not running"
    return 1
  else
    kill ${pid}
  fi
}

freefsm_wait()
{
  local n=60
  echo "Waiting for FreeFSM shutdown:"
  while :
  do
    printf '.'
    pid=$(check_pidfile "${pidfile}" "${procname}")
    if [ -z "${pid}" ]
    then
      printf '\n'
      break
    fi
    sleep 1
    n=$((${n} - 1))
    if [ ${n} -eq 0 -a -f "${pidfile}" ]
    then
      printf "\nForce shutdown"
      kill -9 $(cat "${pidfile}")
      for n in 1 2 3
      do
        printf '.'
        sleep 1
      done
      printf '\n'
      break
    fi
  done
  rm -f "${pidfile}"
  echo "Shutdown complete"
}

run_rc_command "$1"