# deploy/freebsd/freefsm
#!/bin/sh

# PROVIDE: freefsm
# REQUIRE: LOGIN postgresql
# KEYWORD: shutdown

. /etc/rc.subr

name="freefsm"
rcvar="freefsm_enable"

command="/usr/sbin/daemon"
procname="/usr/local/bin/freefsm"
pidfile="/var/run/${name}/${name}.pid"
required_files="${procname}"
start_precmd="${name}_prestart"

freefsm_config="/usr/local/etc/freefsm.conf"

freefsm_prestart()
{
    if [ -f "${freefsm_config}" ]; then
        . "${freefsm_config}"
        export DATABASE_URL SESSION_SECRET PORT ENV STATIC_PATH SETUP_TOKEN
    fi

    if [ ! -f /var/log/freefsm.log ]; then
        touch /var/log/freefsm.log
        chown freefsm:freefsm /var/log/freefsm.log
        chmod 644 /var/log/freefsm.log
    fi

    if [ ! -d "/var/run/${name}" ]; then
        mkdir -p "/var/run/${name}"
    fi
    chown ${freefsm_user}:${freefsm_user} "/var/run/${name}"
    chmod 755 "/var/run/${name}"

    command_args="-p ${pidfile} -u ${freefsm_user} ${procname} -log /var/log/freefsm.log"
}

load_rc_config $name

: ${freefsm_enable:="NO"}
: ${freefsm_user:="freefsm"}

run_rc_command "$1"
