# deploy/freebsd/freefsm
#!/bin/sh

# PROVIDE: freefsm
# REQUIRE: LOGIN postgresql
# KEYWORD: shutdown

. /etc/rc.subr

name="freefsm"
rcvar="freefsm_enable"

load_rc_config $name

: ${freefsm_enable:="NO"}
: ${freefsm_user:="freefsm"}
: ${freefsm_config:="/usr/local/etc/freefsm.conf"}

command="/usr/local/bin/freefsm"
pidfile="/var/run/${name}.pid"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

freefsm_start()
{
    if [ -f "${freefsm_config}" ]; then
        . "${freefsm_config}"
        export DATABASE_URL SESSION_SECRET PORT ENV STATIC_PATH SETUP_TOKEN
    fi

    # Create log file if it doesn't exist
    if [ ! -f /var/log/freefsm.log ]; then
        touch /var/log/freefsm.log
        chown ${freefsm_user}:${freefsm_user} /var/log/freefsm.log
        chmod 644 /var/log/freefsm.log
    fi

    echo "Starting ${name}."
    /usr/sbin/daemon -p ${pidfile} -u ${freefsm_user} ${command} >> /var/log/freefsm.log 2>&1
}

freefsm_stop()
{
    if [ -f "${pidfile}" ]; then
        echo "Stopping ${name}."
        kill -TERM $(cat ${pidfile})
        rm -f ${pidfile}
    else
        echo "${name} is not running."
    fi
}

freefsm_status()
{
    if [ -f "${pidfile}" ]; then
        if kill -0 $(cat ${pidfile}) 2>/dev/null; then
            echo "${name} is running as pid $(cat ${pidfile})."
        else
            echo "${name} is not running (stale pidfile)."
        fi
    else
        echo "${name} is not running."
    fi
}

run_rc_command "$1"
